# drillip.container
# /etc/containers/systemd/drillip.container
[Unit]
Description=Drillip (Sentry SDK compatible error receiver)

[Container]
Image=ghcr.io/philhem/drillip:latest
ContainerName=drillip
Pod=observability.pod
Volume=drillip-data:/data:Z

Environment=DRILLIP_DB=/data/errors.db
Environment=DRILLIP_ADDR=0.0.0.0:8300

PodmanArgs=--memory=16m --memory-swap=16m --cpus=0.25
PodmanArgs=--read-only
PodmanArgs=--security-opt=no-new-privileges:true
PodmanArgs=--cap-drop=ALL

LogDriver=journald

HealthCmd=/drillip health
HealthStartPeriod=5s
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

[Service]
Restart=always
RestartSec=5
