# error-sink.container
# /etc/containers/systemd/error-sink.container
[Unit]
Description=Error sink (Sentry SDK compatible)

[Container]
Image=ghcr.io/philhem/error-sink:latest
ContainerName=error-sink
Pod=observability.pod
Volume=error-sink-data:/data:Z

Environment=ERROR_SINK_DB=/data/errors.db
Environment=ERROR_SINK_ADDR=0.0.0.0:8300

PodmanArgs=--memory=16m --memory-swap=16m --cpus=0.25
PodmanArgs=--read-only
PodmanArgs=--security-opt=no-new-privileges:true
PodmanArgs=--cap-drop=ALL

LogDriver=journald

HealthCmd=/error-sink health
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

[Service]
Restart=always
RestartSec=5
